# files/openclaw/mac/ai.openclaw.macos.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.macos</string>
    <key>ProgramArguments</key>
    <array>
      <string>$HOME/.openclaw/ai.openclaw.macos.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>$HOME</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>$PATH</string>
      <key>NODE_PATH</key>
      <string>$NODE_PATH</string>
      <key>NODE_OPTIONS</key>
      <string>$NODE_OPTIONS</string>
      <key>NODE_EXTRA_CA_CERTS</key>
      <string>$NODE_EXTRA_CA_CERTS</string>
    </dict>
    <key>StandardOutPath</key>
    <string>$HOME/.openclaw/logs/ai.openclaw.macos-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>$HOME/.openclaw/logs/ai.openclaw.macos-stderr.log</string>
  </dict>
</plist>